# config/csmcl-space.service
[Unit]
Description=Cosmic Nexus Web Application
After=network.target mongodb.service

[Service]
Type=simple
User=nginx
Group=nginx
WorkingDirectory=/var/www/cosmic-nexus/backend
ExecStart=/usr/bin/node src/index.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=3000

# Add required capabilities for user management
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_CHOWN
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETUID CAP_SETGID CAP_CHOWN

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home /var/spool/mail

[Install]
WantedBy=multi-user.target
